# appliances/opennebula/etc/rc.local
#!/bin/sh
#
# This script will be executed *after* all the other init scripts.
# You can put your own initialization stuff in here if you don't
# want to do the full Sys V style init stuff.

touch /var/lock/subsys/local

# OpenNebula
ifconfig eth0 promisc
ifconfig tap0 promisc
ifconfig tap1 promisc
/root/bin/datastores-mount
